# last-exercise/master/files/loki/loki.service
[Service]
ExecStart=/usr/bin/loki --config.file=/etc/loki/loki-configuration.yaml


[Install]
WantedBy=multi-user.target